# backend/bot/habit-tracker-worker.service
[Unit]
Description=Habit Tracker Notification Worker
After=network.target

[Service]
User=root
WorkingDirectory=/path/to/your/project/backend
ExecStart=/root/habits_tracker/backend/venv/bin/python bot/notification_worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
